# roles/web/templates/edunotes.service.j2
[Unit]
Description=EduNotes App
After=network-online.target nginx.service
Wants=network-online.target
Requires=nginx.service

[Service]
Type=simple
User=edunotes
Group=edunotes
EnvironmentFile=-/etc/edunotes/env       # ← env が無くても起動継続
WorkingDirectory=/opt/edunotes
ExecStartPre=/usr/bin/install -d -o edunotes -g edunotes -m 0755 /var/edunotes/uploads
ExecStart=/usr/bin/env node server.js     # ← フルパスより柔軟に
Restart=always
RestartSec=3
LimitNOFILE=65535

# npm / Node 用環境変数を安全に継承
Environment=HOME=/home/edunotes
Environment=NODE_ENV=production
Environment=PORT=3000

# Nginx 起動後、Node 側が確実に立ち上がるまでの猶予（数秒）
ExecStartPost=/bin/sleep 3

[Install]
WantedBy=multi-user.target
