# roles/web/templates/edunotes.service.j2
[Unit]
Description=EduNotes App
After=network-online.target nginx.service
Wants=network-online.target
Requires=nginx.service

[Service]
Type=simple
User=edunotes
Group=edunotes
EnvironmentFile=-/etc/edunotes/env
WorkingDirectory=/opt/edunotes

# 事前にアップロード用ディレクトリを用意
ExecStartPre=/usr/bin/install -d -o edunotes -g edunotes -m 0755 /var/edunotes/uploads

# Node アプリ起動
ExecStart=/usr/bin/env node server.js
Restart=always
RestartSec=3

# open files 上限は1回だけ明示
LimitNOFILE=65536

# 実行環境
Environment=HOME=/home/edunotes
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=APP_PORT=3000

ExecStartPost=/bin/sleep 3

# ログ: /var/log/edunotes/app.log へ追記（systemd がディレクトリ自動作成）
LogsDirectory=edunotes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
